# sysconfig_files/etc/systemd/system/fancontrol-timer.service
[Unit]
Description=Runs the fancontrol python script

[Service]
Type=simple
ExecStart=/home/chepaz/projekte/raspberry/fancy-fanpi/ffpi_worker.py
